# code/localModerationNSFWJS/local-moderation-nsfwjs.service
[Unit]
Description=Local Moderation NSFWJS Service (Docker)
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=5
WorkingDirectory=/home/tarun-envid/envid-metadata/code/localModerationNSFWJS
ExecStartPre=/usr/bin/docker rm -f local-moderation-nsfwjs
ExecStartPre=/usr/bin/docker build -t local-moderation-nsfwjs:dev .
ExecStart=/usr/bin/docker run --name local-moderation-nsfwjs -p 5082:5082 -e PORT=5082 local-moderation-nsfwjs:dev
ExecStop=/usr/bin/docker rm -f local-moderation-nsfwjs

[Install]
WantedBy=multi-user.target
